# dot_config/autostart/private_aw-qt.desktop
[Desktop Entry]
Exec=/home/chrio/activitywatch/aw-qt
Icon=
Name=aw-qt
Path=
Terminal=True
Type=Application
